# usr/etc/inittab
::sysinit:/usr/etc/init.d/rcS
::askfirst:/userfs/bin/cli
